# airport_upstart.conf
description "airport"
start on startup
task
exec /home/airport/airport/airport_upstarts.sh